# root/etc/services.d/pm2/run
#!/usr/bin/with-contenv bash
echo "starting pm2"
exec /usr/local/bin/pm2 start /config/apps/process.json --no-daemon